# deploy/systemd/polymarket-paper-bot.service
[Unit]
Description=Polymarket Paper Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/obsidian-hermes-vault/projects/polymarket-5min-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/obsidian-hermes-vault/projects/polymarket-5min-bot/cli.py run --mode paper --runtime-dir /root/obsidian-hermes-vault/projects/polymarket-5min-bot/data/runtime --sleep-seconds 15
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
